# roles/i3/templates/betterlockscreen@.service
[Service]
ExecStart=
ExecStart=/home/{{ user.name }}/.local/share/scripts/lockscreen/lockscreen.sh lock